# backend/scheduler/crontab
*/5 * * * * bin/console app:sources:poll -v
